# examples/a42_ap_withroot/root/etc/crontabs/root
# This script will restart the wifi at 4:00 AM
0 4 * * * /sbin/wifi reload
