# examples/a42_ap_withroot/root/etc/uci-defaults/20-password
PASS="themostsecretpassword"
(echo "$PASS"; sleep 1; echo "$PASS";) | passwd
